# deploy/accredit-api.service
[Unit]
Description=Accredit.vip Backend API
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/accredit.vip/backend
EnvironmentFile=/home/deploy/accredit.vip/backend/.env
ExecStart=/home/deploy/accredit.vip/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
