# deploy/accredit-web.service
[Unit]
Description=Accredit.vip Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/accredit.vip/frontend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /home/deploy/accredit.vip/frontend/node_modules/.bin/next start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
